# volatility_bot.service
[Unit]
Description=Crypto Volatility Bot Service
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/opt/bots/volatility_bot

# Live logs
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONIOENCODING=UTF-8"
Environment="BOT_ENV=prod"

ExecStart=/opt/bots/volatility_bot/venv/bin/python /opt/bots/volatility_bot/main.py

Restart=always
RestartSec=5

# Prevent endless crash-loop hammering
StartLimitIntervalSec=300
StartLimitBurst=10

# Give it time to stop gracefully
TimeoutStopSec=30

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
